# schedules/systemd/gen-checksums.timer
[Unit]
Description=Generate checksums for media files
RefuseManualStart=no
RefuseManualStop=no

[Timer]
OnCalendar=*:0/15
Unit=gen-checksums.service

[Install]
WantedBy=timers.target
