# homedir/.crontab
# Rust
0 1 * * *  rustup update
# Brew
0 2 * * *  brew update
0 3 * * *  brew upgrade
0 4 * * *  brew cleanup
# Cargo
0 5 * * *  cargo install-update -a
# npm
0 6 * * *  npm update npm -g
0 7 * * *  npm update -g
# Zplug
0 8 * * *  zplug update
# Vim
0 9 * * *  vim +PlugUpgrade +PlugUpdate +qa
